# scripts/box/kkernel.service
# kkernel.service: systemd USER unit for the khive substrate daemon.
#
# Install path: ~/.config/systemd/user/kkernel.service (installed by bootstrap.sh).
#
# kkernel mcp --daemon runs khive_runtime::daemon::run_daemon() in the FOREGROUND
# (binds ~/.khive/khived.sock, then blocks in a tokio::select! accept-loop vs.
# SIGTERM/SIGINT). It does not fork or self-daemonize, so Type=simple is correct
# here: systemd's own process supervision IS the daemonization.
#
# kkernel loads ~/.khive/.env itself via load_khive_dotenv() (crates/kkernel/src/main.rs),
# called unconditionally at the top of main() before argument parsing. Real
# environment variables always win over values in that file, and a missing file
# is not an error. No EnvironmentFile= directive is needed in this unit.
#
# For this unit to start on boot and keep running after you log out of the SSH
# session that installed it, enable lingering for this user once:
#   loginctl enable-linger $USER

[Unit]
Description=khive kernel daemon (kkernel mcp --daemon)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.cargo/bin/kkernel mcp --daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
